# contrib/init/potahcoind.conf
description "Potahcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env POTAHCOIND_BIN="/usr/bin/potahcoind"
env POTAHCOIND_USER="potahcoin"
env POTAHCOIND_GROUP="potahcoin"
env POTAHCOIND_PIDDIR="/var/run/potahcoind"
# upstart can't handle variables constructed with other variables
env POTAHCOIND_PIDFILE="/var/run/potahcoind/potahcoind.pid"
env POTAHCOIND_CONFIGFILE="/etc/potahcoin/potahcoin.conf"
env POTAHCOIND_DATADIR="/var/lib/potahcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # potahcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$POTAHCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run potahcoind."
        echo "The setting must appear in $POTAHCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Potahcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$POTAHCOIND_PIDDIR"
    chmod 0755 "$POTAHCOIND_PIDDIR"
    chown $POTAHCOIND_USER:$POTAHCOIND_GROUP "$POTAHCOIND_PIDDIR"
    chown $POTAHCOIND_USER:$POTAHCOIND_GROUP "$POTAHCOIND_CONFIGFILE"
    chmod 0660 "$POTAHCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$POTAHCOIND_PIDFILE" \
    --chuid $POTAHCOIND_USER:$POTAHCOIND_GROUP \
    --exec "$POTAHCOIND_BIN" \
    -- \
    -pid="$POTAHCOIND_PIDFILE" \
    -conf="$POTAHCOIND_CONFIGFILE" \
    -datadir="$POTAHCOIND_DATADIR" \
    -disablewallet \
    -daemon

